# etc/cron.daily/z99-reboot
#!/bin/sh

# write out last reboot date to file, then reboot, simples
echo `date +"%Y-%m-%d %H:%M"` > /etc/lastreboot
reboot
