# etc/cron.daily/z99-reboot
#!/bin/sh

# write out last reboot date to file, then reboot, simples
echo `date +"%d-%m-%Y"` > /var/www/lastreboot
reboot
